# sbac-iaip-rpm-installer/iaip-jenkins-rpm/src/main/scripts/postinstall
# echo "rpm postinstall $1"

fail() {
    echo "ERROR: $1"
    exit 1
}

## Get hostname to be substituted in template config files
TMPLHOSTNAME="localhost"
TMPLDOMAINNAME="pacificmetrics.com"

THISHOSTNAME="$(hostname)"
THISDOMAINNAME="$(hostname -d)"

ONLYHOSTNAME=`echo $THISHOSTNAME  |  cut -d. -f 1`

## To derive the database hostname naming convention is flip the "CDE" substring to "DB"
ONLYHOSTNAMEALTEREDFORDB=${ONLYHOSTNAME/cde/db}
ONLYHOSTNAMEALTEREDFORDB=${ONLYHOSTNAMEALTEREDFORDB/CDE/DB}

ONLYHOSTNAMEALTEREDFORDB=${ONLYHOSTNAME/app/db}
ONLYHOSTNAMEALTEREDFORDB=${ONLYHOSTNAMEALTEREDFORDB/APP/DB}

THISHOSTNAMEALTEREDFORDB="${ONLYHOSTNAMEALTEREDFORDB}.${THISDOMAINNAME}"


##echo "Tomcat context.xml configuration"
##echo "Domain:          $THISDOMAINNAME"
echo "Hostname:        $THISHOSTNAME"
##echo "Tmpl Hostname:   $TMPLHOSTNAME"
####echo "Derived DB Host: $THISHOSTNAMEALTEREDFORDB"
##echo "Tmpl Domain: $TMPLDOMAINNAME"


# Set ownership
echo "Updating file ownership: jenkins:jenkins"
/bin/chown -R jenkins:jenkins /opt/jenkins-tomcat
/bin/chown -R jenkins:jenkins /opt/jenkins-tomcathome

# Register jenkinstomcat service
/sbin/chkconfig --add jenkinstomcat || fail "chkconfig failed to register jenkinstomcat"

# Enable jenkinstomcat service at boot-time
/sbin/chkconfig jenkinstomcat --level 345 on || fail "chkconfig failed to enable jenkinstomcat service"

# Start jenkinstomcat service
/sbin/service jenkinstomcat start || fail "failed to start jenkinstomcat service"
